# templates/systemd/supabase-direct-proxy.service
[Unit]
Description=Supabase Direct IPv6 Proxy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/default/supabase-direct-proxy
ExecStart=/usr/local/bin/supabase-direct-proxy
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
